# backend/docker/supervisor/supervisord.conf
; Supervisord gestiona los procesos de larga duración del contenedor (Apache,
; worker de cola y scheduler) como PID 1, reiniciando cualquiera que muera.
; Reemplaza los loops `while true` sueltos del entrypoint, que no se reiniciaban
; si el proceso de fondo se caía → las importaciones quedaban "en cola" para siempre.

[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true
user=root
logfile=/var/www/html/storage/logs/supervisord.log
logfile_maxbytes=10MB
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; --- Apache (servidor web) ---
[program:apache]
command=apache2-foreground
autostart=true
autorestart=true
startretries=20
; Logs de Apache al stdout/stderr del contenedor (visible en `docker logs`).
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; --- Worker de cola ---
; Procesa importaciones de catálogo y conversiones de imágenes en segundo plano.
; --timeout supera el timeout del job de importación (1800s). --max-time recicla
; el proceso cada hora para evitar fugas de memoria; supervisord lo relanza solo.
[program:queue-worker]
command=php /var/www/html/artisan queue:work --sleep=2 --tries=1 --timeout=1860 --max-time=3600 --no-interaction
directory=/var/www/html
user=www-data
environment=HOME="/tmp"
numprocs=1
autostart=true
autorestart=true
startsecs=5
startretries=20
; Cierre ordenado: deja terminar el job en curso antes de matar el worker.
stopsignal=TERM
stopwaitsecs=1860
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
priority=20

; --- Vite dev server (solo desarrollo) ---
; Sirve los assets de la SPA con HMR en :5173. Antes había que lanzarlo a mano
; (`npm run dev`) tras cada arranque del contenedor; si no, la página quedaba en
; negro porque `public/hot` apuntaba a un servidor caído.
; En producción (Render) los assets se precompilan con `npm run build`, así que
; el guard sale con exit 0 y supervisord deja el programa en EXITED sin
; reiniciarlo (autorestart=unexpected). Corre como root para poder escribir
; node_modules/.vite-cache y public/hot bajo el bind mount.
[program:vite]
command=bash -c 'if [ -n "${RENDER:-}" ] || [ "${APP_ENV:-}" = "production" ]; then echo "vite deshabilitado en producción (assets precompilados)"; exit 0; fi; cd /var/www/html && exec npm run dev'
directory=/var/www/html
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=0
startretries=20
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/vite.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
priority=15

; --- Scheduler ---
; Loop largo (el comando no termina) para que supervisord lo trate como proceso
; supervisado; si el loop muere, lo reinicia.
[program:scheduler]
command=bash -c "while true; do php /var/www/html/artisan schedule:run --no-interaction; sleep 60; done"
directory=/var/www/html
user=www-data
environment=HOME="/tmp"
numprocs=1
autostart=true
autorestart=true
startretries=20
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
priority=30
